# packages/addons/service/fbsplash/source/system.d/service.fbsplash.service
[Unit]
Description=Framebuffer Splash Service
After=graphical.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.fbsplash/bin/showsplash"

[Install]
WantedBy=kodi.target
